# packaging/helper/linux/nomicous-inference-helper.service
[Unit]
Description=Nomicous Inference Helper
After=network.target

[Service]
Type=simple
ExecStart=__INSTALL_DIR__/nomicous-inference-helper --no-tray
WorkingDirectory=__INSTALL_DIR__
Restart=on-failure
Environment=HF_CACHE_ROOT=__HOME__/.nomicous/hf/cache
Environment=HELPER_REGISTRY_URL=__REGISTRY_URL__
Environment=HELPER_CORS_ORIGINS=__CORS_ORIGINS__
Environment=INFERENCE_MAX_REQUEST_BODY_BYTES=167772160
Environment=INFERENCE_MAX_ENCODED_IMAGE_BYTES=167772160
Environment=INFERENCE_MAX_DECODED_IMAGE_BYTES=104857600
Environment=INFERENCE_MAX_IMAGE_PIXELS=200000000
Environment=INFERENCE_MAX_TRANSCRIBE_LINES=10000

[Install]
WantedBy=default.target
